# docker/crontab
*/15 * * * * echo "Launching calculatePeak command..."
*/15 * * * * /app/bin/console bot:cron:calculatePeak
*/1 * * * * echo "Ping"
