# webapi/docker/crontab
0 0/4 * * * /app/bin/console bot:cron:calculatePeak